# deploy/systemd/social-content-platform.service
[Unit]
Description=Social Content Platform Console
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/xiao/backend
Environment=APP_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8000
Environment=OWNER_ACCESS_TOKEN=replace-me
Environment=APP_DB_PATH=/var/www/xiao/backend/data/social_content.sqlite3
ExecStart=/var/www/xiao/backend/venv/bin/python main.py
Restart=always
RestartSec=3
# Change these to the Linux account that owns the deployment files.
User=deploy
Group=deploy

[Install]
WantedBy=multi-user.target
